# tenant-manager-app.service
[Unit]
Description=Tenant Manager App - Flask Web Application
Documentation=https://github.com/company/tenant-manager-app
After=network.target
Wants=network.target

[Service]
Type=exec
User=tmapp
Group=tmapp
WorkingDirectory=/opt/tenant-manager-app
Environment=PATH=/opt/tenant-manager-app/venv/bin
Environment=FLASK_ENV=production
Environment=PYTHONPATH=/opt/tenant-manager-app

# Production WSGI server configuration
ExecStart=/opt/tenant-manager-app/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 4 \
    --worker-class sync \
    --worker-connections 1000 \
    --max-requests 1000 \
    --max-requests-jitter 100 \
    --timeout 120 \
    --keep-alive 60 \
    --preload \
    --access-logfile /opt/tenant-manager-app/logs/access.log \
    --error-logfile /opt/tenant-manager-app/logs/error.log \
    --log-level info \
    --capture-output \
    wsgi_production:application

# Graceful reload and restart
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/tenant-manager-app/logs /opt/tenant-manager-app/flask_session
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tenant-manager-app

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
